# systemd/tuxedo-perkey.service
[Unit]
Description=TUXEDO Keyboard RGB Per-Key Driver and Config
After=multi-user.target

[Service]
Type=oneshot
# __TUXEDO_PERKEY_DIR__ is replaced by install.sh with the actual path
Environment="TUXEDO_PERKEY_DIR=__TUXEDO_PERKEY_DIR__"
ExecStart=/bin/bash -c 'modprobe tuxedo_nb04_rgb_perkey 2>/dev/null || insmod ${TUXEDO_PERKEY_DIR}/kernel/tuxedo_nb04_rgb_perkey.ko'
ExecStartPost=/bin/bash -c '${TUXEDO_PERKEY_DIR}/app/.venv/bin/python ${TUXEDO_PERKEY_DIR}/app/apply_config.py'
User=root
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
